# lab2/backup_daemon.service
[Unit]
Description=Backup

[Service]
ExecStart=/usr/bin/backup_daemon
StandardOutput=syslog
Restart=always

[Install]
WantedBy=default.target
